# modules/utilities/unix/zeroclaw/files/zeroclaw.service
[Unit]
Description=ZeroClaw AI Agent for SecGen Hackerbot
Documentation=https://github.com/zeroclaw-labs/zeroclaw
After=network.target inspircd.service ollama.service
Wants=inspircd.service ollama.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/zeroclaw
ExecStart=/opt/zeroclaw/zeroclaw --config /opt/zeroclaw/config.toml
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=zeroclaw

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/opt/zeroclaw/workspace /var/lib/zeroclaw /var/log/zeroclaw

# Environment
Environment="RUST_LOG=info"
Environment="ZEROCLOW_WORKSPACE=/opt/zeroclaw/workspace"

[Install]
WantedBy=multi-user.target
